# conf/app.conf
[program:app]
command=java  -jar  /app/app-0.0.1-SNAPSHOT.jar 
stdout_logfile =/tmp/app.log